# cheggwpt/languagetool/container_confs/etc/supervisor.d/languagetool.conf

[program:languagetool]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=/usr/bin/java -cp /usr/local/languagetool/languagetool-server.jar org.languagetool.server.HTTPServer --port 8081 --public -v --allow-origin \"*\"
numprocs=1
autostart=true
autorestart=true
